# rc.local
#!/bin/sh -e
export HOME=/root

sleep 1
eval "cd /root/aios && ./run.sh"

exit 0
